# examples/systemd/romm-fuse.service
[Unit]
Description=RomM FUSE Filesystem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /media/fat
ExecStartPre=/bin/chmod 755 /media/fat
ExecStart=/usr/local/bin/romm-fuse /media/fat --profile mister --foreground --allow-other
ExecStop=/usr/bin/fusermount3 -u /media/fat
Restart=on-failure
RestartSec=10

# Environment (alternative to config file)
# Environment=ROMM_URL=http://192.168.1.100:3000
# Environment=ROMM_TOKEN=rmm_your_token_here

[Install]
WantedBy=multi-user.target
